# rom/usr/lib/opkg/info/fp.list
/lib/modules/2.6.36/fp.ko
/etc/init.d/fp.init
