# rom/usr/lib/opkg/info/url_safe.list
/etc/init.d/safe.init
/lib/modules/2.6.36/safe.ko
