# rom/usr/lib/opkg/info/url_log.list
/lib/modules/2.6.36/tracert.ko
/etc/init.d/tracert.init
