# rom/usr/lib/opkg/info/web_auth.list
/etc/init.d/auth.init
/lib/modules/2.6.36/auth.ko
