# usr/lib/opkg/info/cloud.list
/etc/init.d/cloud
/bin/cloud_client
/bin/cloud_exchange
/bin/tspeed
