# usr/lib/opkg/info/ioos.list
/etc/config/qos_rule
/bin/setsysflag
/bin/resetd
/etc/init.d/ioos
/bin/ioos
/etc/firmware
