# usr/lib/opkg/info/l7.list
/bin/gdbserver
/etc/init.d/l7
/etc/rconf_check
/etc/rconf.tar.gz
