# usr/lib/opkg/info/mu.list
/etc/init.d/mu
/etc/config/host
/bin/mu
