# usr/lib/opkg/info/sqos.list
/bin/qos.sh
/etc/init.d/sqos
/lib/modules/2.6.36/sqos.ko
